# crontab
# m h  dom mon dow   command
  0 *  *   *   *     jsub -quiet -N forrestbot -mem 500M /data/project/forrestbot/forrestbot/forrestbot.sh
